# flask_moudle-master/etc/supervisord/car_ins_agent_mgr.ini
[program:car_ins_agent_mgr]
directory=/home/toptrain/workspace/car_ins_agent_mgr
command=gunicorn --workers 2 --threads 20 -k gevent --timeout 5 --user toptrain --group toptrain -b 127.0.0.1:36200 manage:app --log-level info
user=toptrain
autostart=true
startsecs=2
startretries=10
autorestart=true
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
redirect_stderr = true
stdout_logfile_maxbytes = 500MB
stdout_logfile_backups = 2
stdout_logfile = /home/toptrain/workspace/car_ins_agent_mgr/log/car_ins_agent_mgr.log

[program:car_ins_agent_mgr-login_keeper]
directory=/home/toptrain/workspace/car_ins_agent_mgr
command=python manage.py login_keeper
user=toptrain
autostart=true
startsecs=2
startretries=10
autorestart=true
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
redirect_stderr = true
stdout_logfile_maxbytes = 500MB
stdout_logfile_backups = 2
stdout_logfile = /home/toptrain/workspace/car_ins_agent_mgr/log/login_keeper.log